# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

# uvicorn (FastAPI). Binds 127.0.0.1:8765 — the CLI refuses non-loopback by
# design. Runs as the unprivileged app user.
[program:jobhunter]
command=jobhunter --no-browser
user=app
priority=10
environment=HOME="/home/app",JOBHUNTER_WEB_PORT="8765"
autorestart=true
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Caddy. Runs as root so it can bind 80/443 in prod. Reads the baked Caddyfile;
# cert/data storage goes to /data (mount a volume there to persist certs).
[program:caddy]
command=caddy run --config /etc/caddy/Caddyfile --adapter caddyfile
user=root
priority=20
environment=XDG_DATA_HOME="/data",XDG_CONFIG_HOME="/config"
autorestart=true
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
